# recipes-common/rdm/files/apps-rdm.service
[Unit]
Description= Apps Download Service
After=swupdate.service

[Service]
SyslogIdentifier=apps-rdm
ExecStart=/bin/sh -c '/usr/bin/apps_rdm.sh'
