# recipes-common/rdm/files/apps_rdm.path
[Unit]
Description= Apps Download Path
OnFailure=path-fail-notifier@%n.service

[Path]
PathExists=/tmp/.xconfssrdownloadurl
Unit=apps-rdm.service

[Install]
WantedBy=multi-user.target
